# deploy/systemd/ai-customer-service-backup.timer
[Unit]
Description=Nightly AI Customer Service backup

[Timer]
OnCalendar=*-*-* 02:30:00
RandomizedDelaySec=15m
Persistent=true
Unit=ai-customer-service-backup.service

[Install]
WantedBy=timers.target
